# src/ToolboxAdminApi-oneclick/deploy/baota/toolbox-admin.service
[Unit]
Description=Toolbox Admin API
After=network.target

[Service]
Type=simple
WorkingDirectory=/www/wwwroot/gjx.vst76.cn
Environment=TOOLBOX_HOST=127.0.0.1
Environment=TOOLBOX_PORT=5088
Environment=TOOLBOX_ADMIN_TOKEN=请改成你的强密码
ExecStart=/usr/bin/python3 /www/wwwroot/gjx.vst76.cn/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
